# systemd/peer-list-check.timer
[Unit]
Description=Run the peer list check every 300 seconds.

[Timer]
OnBootSec=300
OnUnitActiveSec=300
AccuracySec=1ms

[Install]
WantedBy=timers.target
